# packages/bsp/jethub/jn5189/jn5189.service
[Unit]
Description=JN5189
After=amlogicw1.target

[Service]
ExecStart=/usr/lib/firmware/jn5189/jn5189_func.sh start

[Install]
WantedBy=multi-user.target
